# etc/systemd/system/tileserver-rebuild.timer
[Unit]
Description=Yearly tile rebuild

[Timer]
OnCalendar=*-01-01 03:00 Asia/Tokyo
Persistent=true
RandomizedDelaySec=30m

[Install]
WantedBy=timers.target
